# supervisord.conf
[supervisord]
loglevel=debug

[program:stucco-ui]
command=http-server /stucco/ui -p 8000
process_name=stucco-ui_%(process_num)s
numprocs=1
autostart=true
autorestart=true
